# conf/pilea.cron
*/10  *  *  *  * #USER# #DESTDIR#/bin/console pilea:fetch-data false