# conf/pilea.cron
*/10  *  *  *  * #USER# php7.2 #DESTDIR#/bin/console pilea:fetch-data false
